# ops/launchd/io.tabber.anthropic-research-kr.translation-batch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.tabber.anthropic-research-kr.translation-batch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/tabber/Anthropic-Research-KR/scripts/run_translation_batch_loop.sh</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/tabber/Anthropic-Research-KR</string>

  <key>RunAtLoad</key>
  <true/>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>ProcessType</key>
  <string>Background</string>

  <key>StandardOutPath</key>
  <string>/Users/tabber/Anthropic-Research-KR/logs/codex/translation-batch.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/tabber/Anthropic-Research-KR/logs/codex/translation-batch.stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/tabber</string>
    <key>TRANSLATION_BATCH_ENV_FILE</key>
    <string>/Users/tabber/Anthropic-Research-KR/ops/codex/translation-batch.env</string>
  </dict>
</dict>
</plist>
